# c-socket-server/dlt-daemon.service
[Unit]
Description=DLT Daemon Service
After=network.target

[Service]
ExecStart=/usr/local/bin/dlt-daemon
Restart=always
RestartSec=5
User=nobody
Group=nogroup
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target